# cron.d/ldap-acl-syncer
* * * * * root /usr/bin/ldap-acl-syncer.pl
